# config/systemd/beestgraph-watcher.service
[Unit]
Description=beestgraph vault inbox watchdog
Documentation=https://github.com/terbeest/beestgraph
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=turbobeest
Group=turbobeest
WorkingDirectory=/mnt/walnut-drive/dev/beestgraph
EnvironmentFile=-/mnt/walnut-drive/dev/beestgraph/docker/.env

ExecStart=/home/turbobeest/.local/bin/uv run python -m src.pipeline.watcher --config /mnt/walnut-drive/dev/beestgraph/config/beestgraph.yml

# Restart on failure with backoff
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beestgraph-watcher

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/turbobeest/vault /mnt/walnut-drive/dev/beestgraph /home/turbobeest/.cache/uv
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
