# config/systemd/beestgraph-watcher.service
[Unit]
Description=beestgraph vault inbox watchdog
Documentation=https://github.com/terbeest/beestgraph
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=turbobeest
Group=turbobeest
WorkingDirectory=/home/turbobeest/beestgraph
EnvironmentFile=-/home/turbobeest/beestgraph/docker/.env
Environment=PATH=/home/turbobeest/.local/bin:/usr/local/bin:/usr/bin:/bin

# Load secrets from 1Password, then start the watcher
ExecStart=/bin/bash -c 'source <(bash /home/turbobeest/beestgraph/scripts/load-secrets.sh) && exec /home/turbobeest/.local/bin/uv run python -m src.pipeline.watcher --config /home/turbobeest/beestgraph/config/beestgraph.yml'

# Restart on failure with backoff
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beestgraph-watcher

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/turbobeest/vault /home/turbobeest/beestgraph /home/turbobeest/.cache/uv /home/turbobeest/.config/op
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
