# config/systemd/beestgraph-bot.service
[Unit]
Description=beestgraph Telegram bot
Documentation=https://github.com/terbeest/beestgraph
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=turbobeest
Group=turbobeest
WorkingDirectory=/mnt/walnut-drive/dev/beestgraph
EnvironmentFile=-/mnt/walnut-drive/dev/beestgraph/docker/.env

ExecStart=/home/turbobeest/.local/bin/uv run python -m src.bot.telegram_bot

# Restart on failure with backoff
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beestgraph-bot

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/mnt/walnut-drive/dev/beestgraph /home/turbobeest/.cache/uv
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
